# tests/supervisord.conf
; Test supervisor configuration file

[supervisord]
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0
pidfile=/tmp/supervisord.pid
nodaemon=true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:logging]
command = supervisor_logging.py
events = PROCESS_LOG

[program:sleep]
command=sleep 10000
startsecs=3

[program:false]
command=false
startsecs=3
startretries=2

[program:true]
command=true
startsecs=0
