# tests/supervisord.conf
; Test supervisor configuration file

[inet_http_server]
port=0.0.0.0:9001

[supervisord]
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0
pidfile=/tmp/supervisord.pid
nodaemon=true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:logging]
command = ./supervisor_logging/__init__.py
events = PROCESS_LOG

[program:messages]
command=./tests/messages
stdout_events_enabled=1
